# task05/Solution1-rules/etc_udev_rules.d_90-load-module-for-usb-keyboard.rules
ACTION=="add", SUBSYSTEM=="usb", DRIVER=="usbhid", RUN+="/tmp/add.sh"
ACTION=="remove", SUBSYSTEM=="input", RUN+="/tmp/remove.sh"
